# etc/cron.d/rbup-storage1.cron
# File: /etc/cron.d/rbup-storage1.cron
# Package: rbup
05	1	*	*	*	root	/usr/share/rbup/rbup.sh -c /etc/rbup/storage1.conf 1>dev/null 2>&1
